# 94-trackpoint.rules
KERNEL=="serio2", ACTION=="add", RUN+="/bin/echo something > /dev/ttyprintk" 
